# scrapers/crontab
# Author(s): Ethan Ratnofsky
0 0 * * * node scrape.js
